# .config/zsh/profile.d/dircolors
eval $(dircolors $XDG_CONFIG_HOME/zsh/dircolors-solarized/dircolors.256dark)
